# .cloud/docker/supervisor/nginx.conf
[program:nginx]
command = nginx
startsecs = 2

; Redirect output to syslog
stdout_syslog = true
stderr_syslog = true
redirect_stderr = true

; Apache is rather important
priority = 200
